# rc.local
https://www.zhangshengdong.com/post/esxi_dms/
