# rc.local
#!/bin/sh
# https://www.zhangshengdong.com/post/esxi_dms/
mount -t nfs 192.168.10.36:/volume1/nfs /mnt/nfs
ethtool -s eno1 wol g
exit 0
